# examples/systemd/weewx-clearskies-realtime.service
[Unit]
Description=weewx-clearskies-realtime
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
EnvironmentFile=/etc/weewx-clearskies/secrets.env
ExecStart=/home/ubuntu/repos/weewx-clearskies-realtime/.venv/bin/weewx-clearskies-realtime
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
