# examples/systemd/weewx-clearskies-config.service
[Unit]
Description=weewx-clearskies-config UI
After=network.target

[Service]
Type=simple
User=clearskies
Group=clearskies
WorkingDirectory=/etc/weewx-clearskies
ExecStart=/usr/local/bin/weewx-clearskies-config --bind :: --port 9876
# If installed in a venv: ExecStart=/path/to/venv/bin/weewx-clearskies-config --bind :: --port 9876
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
# Required so the bootstrap POST can write the generated secrets.env file.
# ProtectSystem=strict makes the filesystem read-only by default; this
# exemption grants write access to the config directory only.
ReadWritePaths=/etc/weewx-clearskies
PrivateTmp=true

[Install]
WantedBy=multi-user.target
